# GrabAccess_SourceCode/Grab2/docs/grub-boot-success.service
[Unit]
Description=Mark boot as successful

[Service]
Type=oneshot
ExecStart=/usr/sbin/grub2-set-bootflag boot_success
