# GrabAccess_SourceCode/Grab2/docs/grub-boot-success.timer
[Unit]
Description=Mark boot as successful after the user session has run 2 minutes
ConditionUser=!@system

[Timer]
OnActiveSec=2min
